# crontab
47 4 * * * /usr/bin/certbot renew 2>&1 > /var/log/https/certbot-renew.log

0 4 * * * /usr/bin/logrotate /etc/logrotate.d/certbot-renew
